# base/ruby/files/etc/profile.d/ruby.sh
export PATH=$PATH:/usr/local/bundle/bin
